# services.d/engine/run
#!/usr/bin/with-contenv sh
exec python /srv/banditio.engine/engine.py