# services.d/mitmproxy/run
#!/usr/bin/with-contenv sh
exec mitmdump --conf /srv/banditio.engine/mitmproxy/mitmdump.conf